# production/jobanalyzer-server/jobanalyzer.cron
# Crontab to run on the jobanalyzer server.

MAILTO=larstha@uio.no

@reboot $HOME/sonar/start-infiltrate.sh
@reboot $HOME/sonar/start-sonalyzed.sh

# UiO ML nodes cluster "mlx.hpc.uio.no"

2,7,12,17,22,27,32,37,42,47,52,57 * * * * $HOME/sonar/scripts/mlx.hpc.uio.no/load-5m-and-upload.sh
2 0-23/2 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/cpuhog.sh
2 0-23/2 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/deadweight.sh
3 0-23 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/load-1h.sh
3 0 1-31 * * $HOME/sonar/scripts/mlx.hpc.uio.no/load-24h.sh
3 0-23/12 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/cpuhog-report.sh
3 0-23/12 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/deadweight-report.sh
4 0-23 * * * $HOME/sonar/scripts/mlx.hpc.uio.no/upload-data.sh

# UiO "Fox" supercomputer "fox.educloud.no"

2,7,12,17,22,27,32,37,42,47,52,57 * * * * $HOME/sonar/scripts/fox.educloud.no/load-5m-and-upload.sh
3 0-23 * * * $HOME/sonar/scripts/fox.educloud.no/load-1h.sh
3 1 1-31 * * $HOME/sonar/scripts/fox.educloud.no/load-24h.sh
4 0-23 * * * $HOME/sonar/scripts/fox.educloud.no/upload-data.sh
2 0-23/2 * * * $HOME/sonar/scripts/fox.educloud.no/deadweight.sh
3 0-23/12 * * * $HOME/sonar/scripts/fox.educloud.no/deadweight-report.sh

# Sigma2 "Saga" supercomputer "saga.sigma2.no"

2,7,12,17,22,27,32,37,42,47,52,57 * * * * $HOME/sonar/scripts/saga.sigma2.no/load-5m-and-upload.sh
3 0-23 * * * $HOME/sonar/scripts/saga.sigma2.no/load-1h.sh
3 1 1-31 * * $HOME/sonar/scripts/saga.sigma2.no/load-24h.sh
4 0-23 * * * $HOME/sonar/scripts/saga.sigma2.no/upload-data.sh
