# production/jobanalyzer-server/sonalyzed.cron
# Crontab to run on the jobanalyzer server to keep sonalyzed running (only).

MAILTO=larstha@uio.no

# TODO: This is pretty primitive and we want somthing better.
@reboot $HOME/sonar/start-sonalyzed.sh
