# production/ml-nodes/sonar-runner.cron
# Cron script to run on compute nodes.  This just runs sonar every five minutes and captures its
# output.  Run this via crontab as the analysis user, eg
#
#   crontab sonar-runner.cron

0-59/5 * * * * $HOME/sonar/sonar.sh
