# production/sonar-nodes/fox.educloud.no/sonar-runner.cron
# Cron file for fox nodes that have slurm (compute and gpu nodes).
#
# Copy this into /etc/cron.d as root:
#
#  # cp sonar-runner.cron /etc/cron.d/sonar
#
# Be sure that the ownership is root.root.

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=larstha@uio.no
0-59/5 * * * * sonar /cluster/var/sonar/bin/sonar.sh
