# production/sonar-nodes/mlx.hpc.uio.no/sonar-runner.cron
# Cron script to run on compute nodes.  This just runs sonar every five minutes and captures its
# output.  Copy this file to /etc/cron.d/sonar.

MAILTO=larstha@uio.no

0-59/5 * * * * sonar-runner $HOME/sonar.sh
